# bin/decrypt.service
[Unit]
Description=Decryption

[Service]
ExecStart=/home/client/Documents/GSLW_TP3_stationcryptage/TP3_stationcryptage/TP3_stationcryptage/bin/decrypt.sh

[Install]
WantedBy=multi-user.target
